# mac-agent/com.cura-stro.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cura-stro.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>__SCRIPT_DIR__/.venv/bin/python3</string>
        <string>__SCRIPT_DIR__/agent.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__SCRIPT_DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>AGENT_PORT</key>
        <string>7777</string>
        <key>AGENT_TOKEN</key>
        <string>__AGENT_TOKEN__</string>
        <key>PIXINSIGHT_BIN</key>
        <string>/Applications/PixInsight/PixInsight.app/Contents/MacOS/PixInsight</string>
        <key>BATCH_SCRIPT</key>
        <string>__SCRIPT_DIR__/cura_batch.js</string>
        <key>WBPP_SCRIPT</key>
        <string>/Applications/PixInsight/src/scripts/BatchProcessing/WeightedBatchPreProcessing.js</string>
        <key>WORK_DIR</key>
        <string>/tmp/cura-stro-jobs</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/cura-stro-agent.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/cura-stro-agent.err</string>
</dict>
</plist>
